# templates/singularity.service.erb
# Autogenerated by Puppet from singularity/templates/singularity.service.erb
# Do not edit

[Unit]
Description=Singularity
<%if ! [nil, :undefined, :undef].include?(scope.lookupvar("singularity::systemd_unit_want")) -%>
Wants=<%= scope.lookupvar("singularity::systemd_unit_want")%>
<% end -%>
<%if ! [nil, :undefined, :undef].include?(scope.lookupvar("singularity::systemd_unit_after")) -%>
After=<%= scope.lookupvar("singularity::systemd_unit_after") %>
<% end -%>

[Service]
ExecStart=/bin/bash -c "set -o errexit -o nounset; exec <%=@install_dir%>/run.sh >><%=@log_directory%>/singularity.log 2>&1"

Restart=always
RemainAfterExit=no
SyslogIdentifier=singularity
User=<%= scope.lookupvar("singularity::user") %>
Group=<%= scope.lookupvar("singularity::group") %>

[Install]
WantedBy=multi-user.target
